# systemd/sisa-results.timer
[Unit]
Description=Check sisa results every 10 minutes
[Timer]
OnCalendar=*:0/10
Unit=sisa-results.service
[Install]
WantedBy=multi-user.target